# deploy/demo/systemd/lattice-demo-reset.service
# Installed by demo-bootstrap.sh with __REPO__ substituted. One reset run;
# scheduling lives in lattice-demo-reset.timer.
[Unit]
Description=Lattice hosted demo — nightly world reset (wipe + reseed)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=__REPO__
ExecStart=__REPO__/deploy/demo/demo-reset.sh
# The wipe rebuilds every Refractor durable, so the next start rescans the whole
# graph before the auth plane answers for any new actor — tens of minutes on a
# 2 vCPU box. provision-demo-operator.sh waits that out (PROVISION_TIMEOUT_SECONDS,
# default 45m), so this budget must exceed it or systemd kills the reset partway
# and the demo comes back without its Loupe operator.
TimeoutStartSec=4200
# demo-up.sh leaves the stack's host processes (processor … facet) running in
# the background, lifecycle-managed by the scripts' own pkills; the default
# control-group cleanup would reap them all the moment the oneshot deactivates.
KillMode=process
