# deploy/demo/systemd/lattice-demo-reset.timer
[Unit]
Description=Nightly reset of the Lattice hosted demo world

[Timer]
OnCalendar=*-*-* 09:10:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
